# dynamicdns/dynamicdns.timer
[Unit]
Description=Run the dynamic dns script every hour

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target

