# external_conn/checkextconn.timer
[Unit]
Description=Run the external connectivity script every hour

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target

